# etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:/sbin/getty -L ttyPS0 115200 vt100
::ctrlaltdel:/bin/umount -a -r
